# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/beetrack-web schedule:run >> /dev/null 2>&1
